# src/influx/cronjob.cron
0 1 * * * /var/backup/influx/backup.sh
0 1 * * * /var/log/influx/logs.sh
